# supervisord-tomcat.conf
[program:tomcat]
command=/opt/tomcat-supervisor.sh
